# config/crontab
6,21,36,51 * * * * /<PATH TO PYTHON HOME>/bin/python /<[PROJECT PATH>/es-gdelt/gdelt_realtime_downloader.py >> ~/logs/cron_gdelt_downloader.log 2>&1
